# bc_vt_desktop/Singularityfiles/services.conf
[program:Xvnc]
command=Xvncstart.sh -d "%(ENV_DISPLAYNAME)s"
autostart=false
killasgroup=true
stopwaitsecs=5
stopasgroup=true
priority=10

[program:websockify]
command=/usr/local/bin/websockify --portfd=/TMPDIR/.novncport --web=/opt/noVNC-1.2.0/ --auth-plugin=RemoteUserAuth --unix-target=/TMPDIR/.vncsocket --file_only 0 
autostart=false
killasgroup=true
priority=20

[program:xfce]
command=xfcestart.sh
autostart=false
killasgroup=true
stopwaitsecs=5
stopasgroup=true
priority=30
